# .miniconda3/etc/profile.d/conda.csh
setenv CONDA_EXE "/Users/diveshjoshi_25/Desktop/Data/dys_ai_bot/backend/.miniconda3/bin/conda";
setenv _CONDA_EXE "/Users/diveshjoshi_25/Desktop/Data/dys_ai_bot/backend/.miniconda3/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/Users/diveshjoshi_25/Desktop/Data/dys_ai_bot/backend/.miniconda3/bin/python";
setenv _CONDA_ROOT "/Users/diveshjoshi_25/Desktop/Data/dys_ai_bot/backend/.miniconda3";

source "/Users/diveshjoshi_25/Desktop/Data/dys_ai_bot/backend/.miniconda3/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"